# deploy/systemd/bitkub-streamlit.service
[Unit]
Description=Bitkub Streamlit Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bitkub
WorkingDirectory=/opt/bitkub
Environment=BITKUB_APP_ROOT=/opt/bitkub
Environment=BITKUB_VENV_PATH=/opt/bitkub/.venv
Environment=BITKUB_CONFIG_PATH=/opt/bitkub/config.json
Environment=BITKUB_DB_PATH=/opt/bitkub/data/bitkub.db
Environment=BITKUB_RUNTIME_STATE_PATH=/opt/bitkub/runtime_state.json
Environment=STREAMLIT_SERVER_PORT=8501
Environment=STREAMLIT_SERVER_ADDRESS=0.0.0.0
EnvironmentFile=-/opt/bitkub/.env
ExecStart=/opt/bitkub/scripts/start_streamlit.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
