# deploy/systemd/bitkub-streamlit.service
[Unit]
Description=Bitkub Bot Streamlit UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bitkub
WorkingDirectory=/opt/bitkub/Bitkub_Bot
Environment=BITKUB_APP_ROOT=/opt/bitkub/Bitkub_Bot
Environment=BITKUB_VENV_PATH=/opt/bitkub/Bitkub_Bot/.venv
Environment=STREAMLIT_SERVER_ADDRESS=0.0.0.0
Environment=STREAMLIT_SERVER_PORT=8501
EnvironmentFile=-/opt/bitkub/Bitkub_Bot/.env
ExecStart=/opt/bitkub/Bitkub_Bot/scripts/start_streamlit.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target