# deploy/systemd/bitkub-engine.service
[Unit]
Description=Bitkub Bot Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bitkub
WorkingDirectory=/opt/bitkub
Environment=BITKUB_APP_ROOT=/opt/bitkub
Environment=BITKUB_VENV_PATH=/opt/bitkub/.venv
Environment=BITKUB_CONFIG_PATH=/opt/bitkub/config.json
Environment=BITKUB_DB_PATH=/opt/bitkub/data/bitkub.db
Environment=BITKUB_RUNTIME_STATE_PATH=/opt/bitkub/runtime_state.json
EnvironmentFile=-/opt/bitkub/.env
ExecStart=/opt/bitkub/scripts/start_engine.sh
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
