# deploy/archive/systemd/bitkub-engine.service
[Unit]
Description=Bitkub Bot Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bitkub
WorkingDirectory=/opt/bitkub/Bitkub_Bot
Environment=BITKUB_APP_ROOT=/opt/bitkub/Bitkub_Bot
Environment=BITKUB_VENV_PATH=/opt/bitkub/Bitkub_Bot/.venv
EnvironmentFile=-/opt/bitkub/Bitkub_Bot/.env
ExecStart=/opt/bitkub/Bitkub_Bot/scripts/start_engine.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target